# rootfs/etc/services.d/honeygain-autoclaim/run
#!/command/with-contenv bash
exec /bin/honeygain_autoclaim_service 
sleep 5