# rootfs/etc/services.d/earnapp-log2ram/run
#!/command/with-contenv bash
exec /bin/earnapp_log2ram_service 
sleep 5